# server/deploy/ueat-embedding-backfill.service
[Unit]
Description=ueat AI embedding backfill
After=network.target postgresql.service ollama.service

[Service]
Type=oneshot
WorkingDirectory=/opt/ueat/server
Environment=NODE_ENV=production
Environment=DATABASE_URL=postgresql:///ueat?host=/var/run/postgresql
Environment=AI_EMBEDDING_PROVIDER=ollama
Environment=OLLAMA_BASE_URL=http://127.0.0.1:11434
Environment=OLLAMA_EMBEDDING_MODEL=bge-m3
Environment=AI_EMBEDDING_VECTOR_DIMENSIONS=1024
Environment=AI_EMBEDDING_TIMEOUT_MS=15000
Environment=AI_PGVECTOR_ENABLED=true
ExecStart=/usr/bin/node dist/tools/embeddingBackfill.js 200
User=root
